# pa_agent/build/com.clipbutler.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS LaunchAgent for CLPBTLR.
  Installs to: ~/Library/LaunchAgents/com.clpbtlr.agent.plist
  Load:   launchctl load ~/Library/LaunchAgents/com.clpbtlr.agent.plist
  Unload: launchctl unload ~/Library/LaunchAgents/com.clpbtlr.agent.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.clpbtlr.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Applications/CLPBTLR.app/Contents/MacOS/CLPBTLR</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/Shared/CLPBTLR/clpbtlr.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/Shared/CLPBTLR/clpbtlr.err</string>

  <key>WorkingDirectory</key>
  <string>/Applications/CLPBTLR.app/Contents/MacOS</string>
</dict>
</plist>
